# debian/postinst
#!/bin/sh

#DEBHELPER#

# FIXME: disabling system nginx for the moment
invoke-rc.d nginx stop || true
update-rc.d nginx disable || true